# sample.monitor.cron
*/10 * * * * root cd {src} && ./monitor.sh
