# sample.nginx.cron
0 0 1,6,11,16,21,26 * * root docker stop nginx && certbot renew && {path-nginx}/copy_all_certs.sh && docker start nginx
