# sample.monitor.cron
* * * * * root cd {src} && ./monitor.sh
